# templates/etc/init/kibana.conf.erb
description "start and stop kibana"
author "Gareth Rushgrove"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 30

chdir /opt/kibana

script
  su www-data -c 'ruby kibana.rb'
end script
